# launchd/com.geothings.geobingan.fridayreport.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.geothings.geobingan.fridayreport</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/geothingsmacbookair/Documents/GitHub/geoBingAn-pdf-sync-tool/run_friday_report.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>5</integer>
        <key>Hour</key>
        <integer>17</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardErrorPath</key>
    <string>/Users/geothingsmacbookair/Documents/GitHub/geoBingAn-pdf-sync-tool/logs/launchd_fridayreport_err.log</string>
</dict>
</plist>
